# deployment/systemd/rich-erp.service
[Unit]
Description=Rich ERP Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/go-crm/erp/.next/standalone
ExecStart=/usr/bin/node server.js
Restart=always
Environment=NODE_ENV=production
Environment=PORT=42113
Environment=HOSTNAME=0.0.0.0
Environment=NEXT_PUBLIC_API_URL=http://localhost:42111

[Install]
WantedBy=multi-user.target
